# templates/initscript/logstash-init.erb
#!/bin/bash
#### BEGIN INIT INFO
## Provides:          <%= title %>
## Required-Start:    $remote_fs $syslog
## Required-Stop:     $remote_fs $syslog
## Should-Start:      $all
## Should-Stop:       $all
## Default-Start:     2 3 4 5
## Default-Stop:      0 1 6
## Short-Description: Logstash
## Description:       Logstash - A tool for managing events and logs
#### END INIT INFO

exit 0

NEEDS TENDER AND LOVE

NAME=<%= title %>
PATH=/usr/sbin:/usr/bin:/sbin:/bin
LOGSTASH_USER=<%= $logstash::params::user %>
LOGSTASH_GROUP=<%= $logstash::params::group %>
LOGSTASH_PIDFILE=<%= logstash_basedir %>/${NAME}.pid

start() {
  echo -n "Starting ${NAME}: "

  if [ "x$USER" != "x$LOGSTASH_USER" ]; then
    sudo -u ${LOGSTASH_USER} java -jar logstash.jar agent -f ${LOGSTASH_CONFIG} -l ${LOGSTASH_LOGPATH}
  else
    /usr/bin/${NAME} -l /var/log/logstash/${NAME}.log -b -p ${LOGSTASH_PIDFILE}
  fi
  echo "done."
}

stop() {
  if [ ! -f ${LOGSTASH_PIDFILE} ]; then
    pgrep -lf ${NAME} >/dev/null 2>&1
    if [ "${?}" -ne "0" ]; then
      echo "${NAME} is running, but the pidfile is not present."
      exit 1
    else
      echo "${NAME} is not running."
      exit 0
    fi
  else
    pgrep -lf ${NAME} >/dev/null 2>&1
    if [ "${?}" -ne "0" ]; then
      echo "${NAME} is not running, but the pidfile is present"
      echo "removing pidfile: ${LOGSTASH_PIDFILE}"
      rm -f ${LOGSTASH_PIDFILE}
      exit 0
    fi
  fi

  echo -n "Shutting down ${NAME}: "
  LOGSTASH_PID=`cat ${LOGSTASH_PIDFILE}`
  kill ${LOGSTASH_PID} && sleep 1

  ps -p ${LOGSTASH_PID} >/dev/null 2>&1
  if [ "${?}" -eq "0" ]; then
    sleep 2
    ps -p ${LOGSTASH_PID} >/dev/null 2>&1
    if [ "${?}" -eq "0" ]; then
      kill -9 ${LOGSTASH_PID}
    fi
  fi

  rm -f ${LOGSTASH_PIDFILE}

  echo "done."
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    sleep 2
    start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
esac

exit 0
